# cod3r.service
[Unit]
Description=Cod3r code editor

[Service]
ExecStart=/cod3r/cod3r.server.py
