# cod3r.socket
[Socket]
ListenStream=12345

[Install]
WantedBy=sockets.target
